# WorkStatusForNetTrafficWithDjango/gunicorn_aliyun_wsfnt.ini
[program:gunicorn_wsfnt]
command=/home/ae/.pyenv/versions/2.7.10/envs/WorkStatusForNetTrafficWithDjango/bin/gunicorn -c gunicornConfig_aliyun.py WorkStatusForNetTrafficWithDjango.wsgi_pro:application
directory=/home/ae/project/WorkStatusForNetTrafficWithDjango/
user=ae
autostart=true
autorestart=true
redirect_stderr=true

[program:celery]
command=/home/ae/.pyenv/versions/2.7.10/envs/WorkStatusForNetTrafficWithDjango/bin/celery -A tasks worker --loglevel=info --beat
directory=/home/ae/project/WorkStatusForNetTrafficWithDjango/task
user=ae
autostart=true
autorestart=true
redirect_stderr=true